# data/cups.service.in
[Unit]
Description=CUPS Printing Service

[Service]
ExecStart=@sbindir@/cupsd -f

[Install]
Also=cups.socket cups.path
WantedBy=printer.target
