# data/cups.path.in
[Unit]
Description=CUPS Printer Service Spool

[Path]
PathExistsGlob=@CUPS_REQUESTS@/d*

[Install]
WantedBy=multi-user.target
